# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#

# prepare rc
cp -rf radio_tools/vimrc /root/.vimrc
cp -rf radio_tools/bashrc /root/.bashrc

# sta files
cp radio_tools/sta/* /etc/config/ && sync
echo 'cp sta files done'

# iptables not compile
#/etc/init.d/firewall stop && echo ok
#echo 'firewall stop done'

#/etc/init.d/network restart && echo ok
#echo 'sta done'

#/radio_tools/py/radioDaemon.py &

exit 0

